# mint-arch-linux/airootfs/etc/systemd/system/mint-blackarch.service
[Unit]
Description=Mint Arch Linux — BlackArch Repository Bootstrap
After=network-online.target
Wants=network-online.target
ConditionPathExists=!/var/lib/pacman/sync/blackarch.db

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mint-blackarch-setup
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
