# airootfs/etc/systemd/system/mint-autoupdate.service
[Unit]
Description=Mint Arch Linux — GitHub Auto Updater
After=network-online.target
Wants=network-online.target
After=nss-lookup.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mint-autoupdate
RemainAfterExit=no
StandardOutput=journal
StandardError=journal
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
